# scheduler/macos/com.tajp.dotfiles-settings-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tajp.dotfiles-settings-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>__HOME__/.local/bin/dotfiles-sync-settings</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>600</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/.local/state/taj-dotfiles/settings-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.local/state/taj-dotfiles/settings-sync-error.log</string>
</dict>
</plist>
