# scheduler/linux/dotfiles-settings-sync.service
[Unit]
Description=Fast-forward taj-p dotfiles and development settings
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/dotfiles-sync-settings
